# uptime-monitor.service
# NOTE: install.sh generates and installs this unit automatically (with your user,
# Python, and install path filled in) - this file is only a template for a MANUAL
# install. Replace YOUR_USER and /path/to/connection-watchman, then:
#   sudo cp uptime-monitor.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now uptime-monitor
[Unit]
Description=Connection Watchman internet monitor (logging daemon)
# time-sync.target: start after NTP is at least underway, so a boot with a stale clock
# (no RTC battery) is less likely to log checks at wrong timestamps.
After=network.target time-sync.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/connection-watchman
ExecStart=/usr/bin/python3 /path/to/connection-watchman/monitor.py
Restart=always
RestartSec=5
# Optional overrides:
# Environment=UPTIME_INTERVAL=5
# Environment=UPTIME_DB=/path/to/connection-watchman/uptime.db
# If this machine runs a VPN, route the probe OUTSIDE the tunnel so it tests the direct link
# (otherwise it measures the VPN path). Find your VPN's bypass mark via `ip rule show`:
# AmbientCapabilities=CAP_NET_ADMIN
# Environment=UPTIME_FWMARK=0xe1f1   # NordVPN / NordLynx

[Install]
WantedBy=multi-user.target
