# uptime-dashboard.service
# NOTE: install.sh generates and installs this unit automatically (with your user,
# Python, and install path filled in) - this file is only a template for a MANUAL
# install. Replace YOUR_USER and /path/to/connection-watchman, then:
#   sudo cp uptime-dashboard.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now uptime-dashboard
[Unit]
Description=Connection Watchman dashboard (web UI)
After=network.target time-sync.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/connection-watchman
ExecStart=/usr/bin/python3 /path/to/connection-watchman/dashboard.py
Restart=always
RestartSec=5
# Environment=UPTIME_PORT=8080

[Install]
WantedBy=multi-user.target
