# uptime-dashboard.service
[Unit]
Description=Internet uptime dashboard (web UI)
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/pi-uptime-monitor
ExecStart=/usr/bin/python3 /home/pi/pi-uptime-monitor/dashboard.py
Restart=always
RestartSec=5
# Environment=UPTIME_PORT=8080

[Install]
WantedBy=multi-user.target
